# crontab
0 3 * * * /home/prophet/etc/tarsnapper.sh >> /tmp/tarsnapper.log 2>&1
